# arena-coligados.service
[Unit]
Description=Arena Coligados Next.js Application
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/srv/arena-coligados
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=10
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=NEXT_PUBLIC_SITE_URL=https://srv998805.hstgr.cloud
EnvironmentFile=/srv/arena-coligados/.env

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=arena-coligados

[Install]
WantedBy=multi-user.target